# letsencrypt-lighttpd.timer
[Unit]
Description=Monthly renewal on letsencrypt's certificates

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
