# letsencrypt-lighttpd.service
[Unit]
Description=Letsencrypt manual renewal

[Service]
Type=oneshot
ExecStart=/home/nessuno/letsencrypt-lighttpd/renew.sh
